# systemd/mumble-ssl.timer
[Unit]
Description=Run mumble-ssl weekly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
